# system-files/etc/systemd/system/dnsmasq.service.d/recovery.conf
[Unit]
ConditionPathNotExists=/boot/firmware/enable-wifi-client
