# system-files/etc/systemd/system/oakhz-audio-events.service
[Unit]
Description=OaKhz Audio Events Manager
After=pulseaudio.service bluetooth.service
Wants=pulseaudio.service bluetooth.service
DefaultDependencies=no

[Service]
Type=simple
User=oakhz
Group=audio
ExecStartPre=/bin/sh -c 'until pactl info >/dev/null 2>&1; do sleep 0.2; done'
ExecStart=/usr/bin/python3 /usr/local/bin/oakhz-audio-events.py
Restart=always
RestartSec=5
Environment="PULSE_SERVER=unix:/run/pulse/native"
Environment="HOME=/home/oakhz"

[Install]
WantedBy=multi-user.target
