# system-files/etc/systemd/system/oakhz-audio-events.service
[Unit]
Description=OaKhz Audio Events Manager
Wants=pulseaudio.service bluetooth.service
DefaultDependencies=no

[Service]
Type=simple
User=oakhz
Group=audio
ExecStartPre=/bin/sh -c 'until pactl info >/dev/null 2>&1; do sleep 0.2; done'
ExecStart=/usr/bin/python3 /usr/local/bin/oakhz-audio-events.py
Restart=always
RestartSec=5
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="PULSE_RUNTIME_PATH=/run/user/1000/pulse"

[Install]
WantedBy=multi-user.target
